# phpfpm/rootfs/etc/services.d/phpfpm/run
#!/usr/bin/with-contenv bash

if [ ! "$ROLE" == "CRON" ]; then
  /usr/sbin/php-fpm${PHP_VERSION} -F --fpm-config /etc/php/${PHP_VERSION}/fpm/php-fpm.conf
fi
